# root/etc/services.d/healthchecks/run
#!/usr/bin/with-contenv bash

cd /app/healtchecks || exit

exec \
	s6-setuidgid abc ./manage.py runserver 0.0.0.0:8000
