# root/etc/services.d/healthchecks/run
#!/usr/bin/with-contenv bash

cd /app/healthchecks || exit

exec \
	s6-setuidgid abc /usr/sbin/uwsgi uwsgi.ini
