# root/etc/services.d/healthchecks_alerts/run
#!/usr/bin/with-contenv bash

cd /app/healthchecks || exit

exec \
	s6-setuidgid abc ./manage.py sendalerts
